# infra/docker/s6/services/init-hebo/up
with-contenv
/etc/s6-overlay/scripts/init-hebo.sh
