# infra/docker/s6/services/svc-gateway/run
#!/bin/sh
case "${HEBO_MODE:-standalone}" in
  standalone|gateway) exec bun /app/gateway/index.js ;;
esac
exec s6-svc -Od .
