# infra/docker/s6/services/svc-gateway/run
#!/bin/sh
exec /etc/s6-overlay/scripts/svc-run.sh svc-gateway "standalone gateway" bun /app/gateway/index.js
