# infra/docker/s6/services/svc-console/run
#!/bin/sh
case "${HEBO_MODE:-standalone}" in
  standalone) exec bun /app/console/serve.ts ;;
esac
exec s6-svc -Od .
