# infra/docker/s6/services/svc-console/run
#!/bin/sh
exec /etc/s6-overlay/scripts/svc-run.sh svc-console "standalone" bun /app/console/serve.ts
