# infra/docker/s6/services/svc-api/run
#!/bin/sh
exec /etc/s6-overlay/scripts/svc-run.sh svc-api "standalone api" bun /app/api/index.js
