# infra/docker/s6/services/svc-mcp/run
#!/bin/sh
case "${HEBO_MODE:-standalone}" in
  standalone|mcp) exec bun /app/mcp/index.js ;;
esac
exec s6-svc -Od .
