# infra/docker/s6/services/svc-auth/run
#!/bin/sh
exec /etc/s6-overlay/scripts/svc-run.sh svc-auth "standalone auth" bun /app/auth/index.js
