# infra/docker/s6/services/svc-auth/run
#!/bin/sh
case "${HEBO_MODE:-standalone}" in
  standalone|auth) exec bun /app/auth/index.js ;;
esac
exec s6-svc -Od .
